# sci-geosciences/gpsd/files/gpsd.conf
# Copyright (c) 2010 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description   "Run gpsd when the system starts"
author        "chromium-os-dev@chromium.org"

start on stopped boot-complete

script
    GPSD_LOG='/tmp/gpsd.log'
    GPS_DEVICE='/dev/ttyUSB2'
    GPSD_PARAMS='-D 7'

    exec > "${GPSD_LOG}" 2>&1

    retry=0
    while [ ! -e "${GPS_DEVICE}" ]
    do
	if [ "${retry}" = "20" ]; then
	    echo "${GPS_DEVICE} not found"
	    exit 1
	fi
	sleep 1
	retry=$((retry + 1))
    done
    echo "Found ${GPS_DEVICE} after ${retry} seconds"
    echo '$GPS_START' > "${GPS_DEVICE}"
    /usr/sbin/gpsd "${GPSD_PARAMS}" "${GPS_DEVICE}"
end script
