# chromeos-base/whining/files/whining_setup.conf
# Copyright 2016 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Setup whining runtime directory."
author          "chromium-os-dev@chromium.org"

start on started boot-services

script
  # Make the logs directory and give ownership to moblab.
  mkdir -p /var/log/whining
  chown -R moblab:moblab /var/log/whining

  scheduler_file="/usr/local/autotest/suite_scheduler.ini"
  if [ ! -e "$scheduler_file" ]
  then
     touch $scheduler_file
  fi
end script
