# chromeos-base/chromeos-debugd/files/debugd.conf
# Copyright (c) 2012 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Chromium OS debug daemon"
author          "chromium-os-dev@chromium.org"

start on x-started
stop on stopping ui or starting halt or starting reboot
respawn

exec /sbin/debugd

